# tools/launchd/com.ai-auto-trading.mainnet-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-auto-trading.mainnet-loop</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec /path/to/ai-auto-trading-service/tools/launchd/run_mainnet_auto_cycle.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/ai-auto-trading-service</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/path/to/ai-auto-trading-service/data/runtime/launchd/mainnet-loop.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/ai-auto-trading-service/data/runtime/launchd/mainnet-loop.stderr.log</string>
</dict>
</plist>
